# roles/flask-dummy/templates/supervisor.ini
[program:{{ app_name }}]
command=/usr/local/bin/uwsgi
  --socket 127.0.0.1:5000
  --logto {{ source_path }}/log/access.log
  --pythonpath /usr/local/lib/python2.7/dist-packages
  --wsgi-file {{ source_path }}/wsgi.py
  --callable {{ wsgi_callable }}
  --max-requests 1000
  --master
  --processes 1
  --chmod
directory={{ source_path }}
autostart=true
autorestart=true
